# container/supervisord.conf
[supervisord]
nodaemon=true

[program:kadalu-mgr]
command=/usr/sbin/kadalu mgr --workdir=/var/lib/kadalu --logdir=/var/log/kadalu
